# system_files/usr/lib/systemd/user/soltros-pipewire-resume.path
[Unit]
Description=Watch Soltros resume stamp and restart PipeWire on change

[Path]
PathChanged=/run/soltros/last_resume
Unit=soltros-pipewire-resume.service

[Install]
WantedBy=default.target
