# system_files/usr/lib/systemd/user/soltros-pipewire-resume.service
[Unit]
Description=Restart PipeWire after system resume for this user

[Service]
Type=oneshot
ExecStart=systemctl --user restart pipewire.service
